# content/til/2026/04/04/launchd-path-whack-a-mole/com.user.argus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.argus</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bitsbyme/tools/argus-build.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/bitsbyme/projects/code-orchestrator</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/argus.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/argus.log</string>
</dict>
</plist>
